# led_control.service
[Unit]
Description=Odroid LED Controller

[Service]
ExecStart=/usr/bin/led_control.sh

[Install]
WantedBy=multi-user.target
